# helios-progress-monitor.timer
[Unit]
Description=Run Helios Progress Monitor every 15 minutes
Requires=helios-progress-monitor.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true
RandomizedDelaySec=2min

[Install]
WantedBy=timers.target